# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
user=1000
command=/usr/local/bin/service.sh
